# mkosi.extra/usr/lib/systemd/system/gdm.service.d/10-x11-socket-fix.conf
[Service]
ExecStartPre=+/usr/bin/mkdir -p /tmp/.X11-unix
ExecStartPre=+/usr/bin/chmod 1777 /tmp/.X11-unix